# db/crontab
0 0 * * 0 /usr/bin/docker exec -t htb "python3 /bot/dbimport.py" 2>&1 >> /var/log/cronjob.htbdocker.log
