# db/crontab
0 0 * * * python3 /bot/dbimport.py >> /var/log/cron.log 2>&1
